# etc/systemd/system/snapper-timeline.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 06:00:00
