# etc/systemd/system/home-narcizo-github.automount
[Unit]
Description=NFS automount (/homes/narcizo/github on cobra)

[Automount]
Where=/home/narcizo/github
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
