# etc/systemd/system/home-narcizo-Music.automount
[Unit]
Description=NFS automount (/share/Music on elefante)

[Automount]
Where=/home/narcizo/Music
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
