# etc/systemd/system/home-narcizo-Videos.automount
[Unit]
Description=NFS automount (/share/Videos on elefante)

[Automount]
Where=/home/narcizo/Videos
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
