# etc/systemd/system/home-narcizo-Pictures.automount
[Unit]
Description=NFS automount (/share/Pictures on elefante)

[Automount]
Where=/home/narcizo/Pictures
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
