# etc/systemd/system/home-narcizo-bin.automount
[Unit]
Description=NFS automount (/homes/narcizo/bin on elefante)

[Automount]
Where=/home/narcizo/bin
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
