# etc/systemd/system/media-Videos.mount
[Unit]
Description=NFS mount (Videos on elefante)
After=wait-for-ping.service

[Mount]
What=192.168.1.104:/Videos
Where=/media/Videos
Type=nfs4
TimeoutSec=20s
Options=_netdev,auto
LazyUnmount=true
ForceUnmount=true

[Install]
WantedBy=remote-fs.target
WantedBy=multi-user.target
