# etc/systemd/system/mnt-backup.automount
[Unit]
Description=NFS automount (linux-backup on elefante)

[Automount]
Where=/mnt/backup
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
