# etc/systemd/system/mnt-gdrive.automount
[Unit]
Description=NFS automount (Google Drive on elefante)

[Automount]
Where=/mnt/gdrive
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
