# etc/systemd/system/mnt-Music.automount
[Unit]
Description=NFS automount (music on cobra)

[Automount]
Where=/mnt/Music
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
