# etc/systemd/system/mnt-backups.automount
[Unit]
Description=NFS automount (linux-backup on cobra)

[Automount]
Where=/mnt/backups
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
