# etc/systemd/system/mnt-Pictures.automount
[Unit]
Description=NFS automount (pictures on cobra)

[Automount]
Where=/mnt/Pictures
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
