# etc/systemd/system/mnt-Videos.automount
[Unit]
Description=NFS automount (videos on cobra)

[Automount]
Where=/mnt/Videos
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
