# etc/systemd/system/mnt-narcizo.automount
[Unit]
Description=NFS automount (narcizo on cobra)

[Automount]
Where=/mnt/narcizo
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
