# etc/systemd/system/snapshot-weekly.service
[Unit]
Description="Weekly rotating-snapshot utility"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapshot-weekly
